# cariot.service
[Unit]
Description=CarIoT Vehicle Monitoring System
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/cariot
ExecStart=/usr/bin/python3 /home/pi/cariot/main.py
Restart=on-failure
RestartSec=15
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
